# buildroot-external/rootfs-overlay/etc/NetworkManager/dispatcher.d/90-openstick-fallback
#!/bin/sh
# NM Dispatcher — trigger OpenStick auto-connect on connectivity loss
#
# NetworkManager calls dispatcher scripts when connection state changes.
# This script triggers the OpenStick auto-connect when connectivity drops.
#
# Arguments from NM: $1=interface $2=action
# connectivity-change action provides CONNECTIVITY_STATE env var
#

INTERFACE="$1"
ACTION="$2"

case "$ACTION" in
  connectivity-change)
    # CONNECTIVITY_STATE: FULL, LIMITED, PORTAL, NONE
    case "${CONNECTIVITY_STATE:-}" in
      NONE|LIMITED)
        /usr/sbin/ga-openstick-autoconnect &
        ;;
    esac
    ;;
  up|down)
    # Also trigger on interface up/down — catches boot without connectivity
    # (connectivity-change may not fire if device was never "full")
    sleep 15  # Wait for NM to settle after interface change
    CONN=$(nmcli -t -f CONNECTIVITY general 2>/dev/null || echo "none")
    case "$CONN" in
      none|limited)
        /usr/sbin/ga-openstick-autoconnect &
        ;;
    esac
    ;;
esac
