# buildroot-ihost/rootfs-overlay/etc/systemd/system/ga-ethernet-guard.service
[Unit]
Description=Apply Ethernet disable state from ga-env.conf
After=NetworkManager.service hassos-config.service
Wants=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/ga-manage-ethernet apply
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
