# buildroot-external/package/telegraf/telegraf.service
[Unit]
Description=Telegraf metrics agent
After=network-online.target docker.service hassio-supervisor.service netbird.service ga-telemetry-consent.service
Wants=network-online.target
# Only start if user consented to metrics collection (marker written by ga-telemetry-consent.service)
ConditionPathExists=/mnt/data/.ga-consent-metrics

[Service]
Type=simple

# Ensure writable directory for env file
ExecStartPre=/bin/mkdir -p /mnt/data/telegraf

# Safe defaults (EnvironmentFile is read BEFORE ExecStartPre, so first boot needs these)
Environment=GA_ENV=dev DEVICE_UUID=unknown GATEWAY_IP=unknown DEVICE_LABEL=unknown GA_BUILD_TIMESTAMP=unknown

# Load GA environment (dev/prod) - check data partition override first, then rootfs default
ExecStartPre=/bin/sh -c '[ -f /mnt/data/ga-env.conf ] && . /mnt/data/ga-env.conf || . /etc/ga-env.conf 2>/dev/null; echo "GA_ENV=${GA_ENV:-dev}" > /mnt/data/telegraf/env'
# Read device UUID into env file for telegraf.conf to use (UUID regex avoids quote escaping issues)
ExecStartPre=/bin/sh -c 'UUID=$(grep -oE "[0-9a-f]{8}-?[0-9a-f]{4}-?[0-9a-f]{4}-?[0-9a-f]{4}-?[0-9a-f]{12}" /mnt/data/supervisor/homeassistant/.storage/core.uuid 2>/dev/null | head -1); echo "DEVICE_UUID=${UUID:-unknown}" >> /mnt/data/telegraf/env'
# Read device label (QR code, e.g. KIB-SON-00001234) for telemetry tagging
ExecStartPre=/bin/sh -c 'LABEL=$(cat /mnt/data/ga-device-label 2>/dev/null); echo "DEVICE_LABEL=${LABEL:-unknown}" >> /mnt/data/telegraf/env'
# Auto-detect default gateway for connectivity monitoring
ExecStartPre=/bin/sh -c 'GW=$(ip route | grep "^default" | head -1 | awk "{print \$3}"); echo "GATEWAY_IP=${GW:-unknown}" >> /mnt/data/telegraf/env'
# Read OS build version for fleet version tracking
ExecStartPre=/bin/sh -c '. /etc/os-release 2>/dev/null; echo "GA_BUILD_TIMESTAMP=${GA_BUILD_TIMESTAMP:-unknown}" >> /mnt/data/telegraf/env'
# EnvironmentFile overrides defaults above on subsequent starts
EnvironmentFile=-/mnt/data/telegraf/env

# Use override config from data partition if present, otherwise rootfs default.
# To switch to debug: cp /etc/telegraf/telegraf-debug.conf /mnt/data/telegraf/override.conf
# To revert:          rm /mnt/data/telegraf/override.conf
ExecStart=/bin/sh -c 'if [ -f /mnt/data/telegraf/override.conf ]; then exec /usr/bin/telegraf --config /mnt/data/telegraf/override.conf; else exec /usr/bin/telegraf --config /etc/telegraf/telegraf.conf; fi'

Restart=always
RestartSec=5
StartLimitIntervalSec=300
StartLimitBurst=10
Nice=10

[Install]
WantedBy=multi-user.target
